# Server/wipi.service
[Unit]
Description="WiPi Server"

[Service]
ExecStart=/usr/bin/wipi

[Install]
WantedBy=multi-user.target

